# cron/crontab
*/5  *  *  *  *   root *study/reloadOpenTreeRepository



